# packaging/rhel-centos/mage
#/bin/sh

###############

# REDHAT chkconfig header

# chkconfig: - 58 74
# description: mage is the script for starting mage on boot.
### BEGIN INIT INFO
# Provides: node
# Required-Start:    $network $remote_fs $local_fs 
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the mage-forever server.
# Description: SAGE uses node.js, forever, and monitor.
### END INIT INFO

FOREVER_BIN=/usr/bin/forever
PID_FILE=/var/run/mage.pid
LOG_FILE=/var/log/mage.log
APP_DIR=/opt/mage/
APP=mage.js
#USER=sage
USER=root
ERR_LOG=/var/log/mage-error.log

. /etc/rc.d/init.d/functions

export PATH=$PATH:/usr/local/bin/

case "$1" in
    start)
        if [ -f $PID_FILE ]
        then
                echo "$PID_FILE exists, process is already running or crashed"
        else
                echo "Starting MAGE:"
                daemon --user=$USER $FOREVER_BIN start\
                 --pidFile $PID_FILE -l $LOG_FILE -a -d\
                  -e $ERR_LOG $APP_DIR$APP
                echo
        fi
        ;;
    stop)
        if [ ! -f $PID_FILE ]
        then
                echo "$PID_FILE does not exist, process is not running"
        else
                echo "Stopping MAGE:"
        #echo "Killing `cat $PID_FILE`"
        echo "Killing $PID_FILE"
        daemon --user=$USER $FOREVER_BIN stop $APP_DIR$APP
        #kill `cat $PID_FILE`;
        #kill 
        rm $PID_FILE;
                echo "MAGE stopped"
        fi
        ;;

    restart)
        if [ ! -f $PID_FILE ]
        then
            echo "$PID_FILE does not exist, process is not running"
        else
            echo "Restarting MAGE:"
            echo "Killing $APP_DIR$APP"
            daemon --user=$USER $FOREVER_BIN stop $APP_DIR$APP
            rm $PID_FILE;
            daemon --user=$USER $FOREVER_BIN start\
             --pidFile $PID_FILE -l $LOG_FILE -a -d\
              -e $ERR_LOG $APP_DIR$APP
            echo "MAGE restarted."
        fi
        ;;
    status)
        if [ ! -f $PID_FILE ]
            then
            echo "$PID_FILE does not exist, MAGE is not running."
        else
            daemon --user=$USER $FOREVER_BIN list --pidfile $PID_FILE
            echo "MAGE is running."
        fi
        ;;
    *)
        echo "Usage: mage {start|stop|restart|status}"
    ;;
esac